# packages/rhel/gdrcopy.service
[Unit]
Description=GDRCopy service
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/libexec/gdrcopy/gdrcopy start
ExecReload=/bin/bash /usr/libexec/gdrcopy/gdrcopy restart
ExecStop=/bin/bash /usr/libexec/gdrcopy/gdrcopy stop

[Install]
WantedBy=multi-user.target

